# .configs/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/talecrew/bin/uwsgi --ini /talecrew/.configs/uwsgi/uwsgi-app.ini